# ops/systemd/qidra-wallet-sync.service
[Unit]
Description=Qidra wallet deposit sync
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/qidra/qidra.env
ExecStart=/usr/local/bin/qidra-wallet-sync.sh
