# ops/systemd/qidra-monitor.service
[Unit]
Description=Qidra production monitor
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/qidra-monitor.sh
